# packaging/etc/sv/inspeqtor/run
#!/bin/sh
exec /usr/bin/inspeqtor -c /etc/inspeqtor
