# packaging/prerm
#!/bin/sh
set -e

if [ -e /etc/systemd/system/inspeqtor.service ]; then
  systemctl stop inspeqtor
  systemctl disable inspeqtor
fi
