# packaging/prerm
#!/bin/sh
set -e

if [ -e /etc/systemd/system/inspeqtor.service ]; then
  systemctl --no-reload disable inspeqtor || :
  systemctl stop inspeqtor || :
fi
